# dot_config/systemd/user/update-check.service
[Unit]
Description=Check for package updates

[Service]
Type=oneshot
ExecStart=%h/scripts/check-package-updates.sh
